# deploy/systemd/evidenceos-daemon.service
# Copyright (c) 2026 Joseph Verdicchio and EvidenceOS Contributors
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=EvidenceOS verification-kernel daemon
After=network.target

[Service]
Type=simple
User=evidenceos
Group=evidenceos
WorkingDirectory=/var/lib/evidenceos
ExecStart=/usr/local/bin/evidenceos-daemon --listen 0.0.0.0:50051 --etl-path /var/lib/evidenceos/etl.log
Restart=on-failure
RestartSec=2

# Hardening (adjust to your deployment)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/evidenceos

[Install]
WantedBy=multi-user.target
